# Config/etc-init/eloquentd.conf
description "X to X logging service"
author "Nathan Wehr <nathanw@evrichart.com>"

start on runlevel [234]
stop on runlevel [0156]

exec /usr/local/sbin/eloquentd --log-level info
respawn
